# com.procore.voice-weekly-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.procore.voice-weekly-update</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/kyleweatherholtz/AI-workshop/procore-voice-weekly-update/weekly_update.py</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>5</integer>
    <key>Hour</key>
    <integer>16</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/Users/kyleweatherholtz/AI-workshop/procore-voice-weekly-update</string>

  <key>StandardOutPath</key>
  <string>/Users/kyleweatherholtz/Library/Logs/procore-voice-update.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kyleweatherholtz/Library/Logs/procore-voice-update-error.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
